# roles/restic/templates/restic-backup.timer.j2
[Unit]
Description=Timer for Restic Backup

[Timer]
OnCalendar={{ restic_timer_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
